# motion.init-RH.in
#!/bin/bash
#
# Startup script for the Motion Detection System 
#
# chkconfig: - 85 15
# description: Motion Detection System.  It is used to detect \
#              movement based on compare images.
# processname: @PACKAGE_NAME@ 
# pidfile: /var/run/@PACKAGE_NAME@.pid
# config: /etc/@PACKAGE_NAME@.conf

# Source function library.
. /etc/rc.d/init.d/functions



motion=${MOTION-@BIN_PATH@/@PACKAGE_NAME@}
prog=@PACKAGE_NAME@
PIDFILE=/var/run/@PACKAGE_NAME@.pid
RETVAL=0


start() {
	echo -n $"Starting $prog: "
	daemon $motion 
	RETVAL=$?
	echo
	[ $RETVAL = 0 ] && touch /var/lock/subsys/@PACKAGE_NAME@
	return $RETVAL
}
stop() {
	echo -n $"Stopping $prog: "
	killproc $motion
	RETVAL=$?
	echo
	[ $RETVAL = 0 ] && rm -f /var/lock/subsys/@PACKAGE_NAME@
}
stopsafe() {
	echo -n $"Stopping $prog: ( for restarting ) "
	killproc $motion
	RETVAL=$?
	sleep 10s
	echo
	[ $RETVAL = 0 ] && rm -f /var/lock/subsys/@PACKAGE_NAME@
}
reload() {
	echo -n $"Reloading $prog: "
	killproc $motion -HUP
	RETVAL=$?
	echo
}

# See how we were called.
case "$1" in
	start)
		start
		;;
	stop)
		stop
		;;
	status)
		status $motion
		RETVAL=$?
		;;
	restart)
		stopsafe
		start
		;;
	condrestart)
		if [ -f /var/run/@PACKAGE_NAME@.pid ] ; then
			stop
			start
		fi
		;;
	reload)
		reload
		;;
	*)
		echo $"Usage: $prog {start|stop|restart|condrestart|reload|status}"
		exit 1
esac

exit $RETVAL
