# motion.service.in
[Unit]
Description=Motion daemon
After=local-fs.target network.target

[Service]
ExecStart=@BIN_PATH@/motion -n
Type=simple
StandardError=null

[Install]
WantedBy=multi-user.target
