# data/motion-dist.service.in
#
#   This file is part of Motion.
#
#   Motion is free software: you can redistribute it and/or modify
#   it under the terms of the GNU General Public License as published by
#   the Free Software Foundation, either version 3 of the License, or
#   (at your option) any later version.
#
#   Motion is distributed in the hope that it will be useful,
#   but WITHOUT ANY WARRANTY; without even the implied warranty of
#   MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
#   GNU General Public License for more details.
#
#   You should have received a copy of the GNU General Public License
#   along with Motion.  If not, see <https://www.gnu.org/licenses/>.
#
#

[Unit]
Description=Motion - Video motion detection and surveillance
Documentation=man:motion(1)
After=local-fs.target network-online.target
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=main
WatchdogSec=30
ExecStart=@bindir@/motion -n
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

# Camera access - run as root or add user to video group
# Can be overridden with "systemctl edit motion"
# Uncomment and set appropriate user/group:
#User=motion
#Group=video

# Don't restart if unconfigured / misconfigured  e.g. daemon disabled
# in defaults file.  See also /usr/include/sysexits.h or sysexits(3)
RestartPreventExitStatus=78

# To tune restart behaviour, see systemd.unit(5) and use
# "systemctl edit motion" to change the following settings:
#StartLimitBurst=
#StartLimitIntervalSec=
#StartLimitAction=
#FailureAction=

# Security hardening (optional, may break some setups)
# Use "systemctl edit motion" to enable these settings.
# See systemd.exec(5) for details.
#PrivateTmp=true
#ProtectHome=read-only
#ProtectSystem=strict
#ReadWritePaths=@localstatedir@/lib/@PACKAGE@

[Install]
WantedBy=multi-user.target
